# src/etc/supervisor.d/gate.ini
[supervisord]
nodaemon=true

[program:redsocks]
command=/usr/bin/redsocks -c /conf/redsocks.conf

[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon --conf-dir=/conf/dnsmasq.d --conf-file=/conf/dnsmasq.conf

[program:pdnsd]
command=/usr/sbin/pdnsd -c /conf/pdnsd.conf

[program:chinadns]
command=/usr/local/bin/chinadns -c /usr/local/share/chnroute.txt -s 202.101.172.35,127.0.0.1:5302 -p 5301 -v

[program:ssledge]
command=/usr/local/bin/ssledge -f /conf/ssledge.conf -n 5000 -N 10000 -m 512
